# jareks.foreman/templates/upstart.conf
start on runlevel [2345]
stop on runlevel [!2345]

chdir {{ app_path }}/current/
exec su dev -c "foreman start"

